# vagrant_minikube/99-vagrant.rules
ACTION=="add", SUBSYSTEM=="scsi" , ATTR{type}=="0|7|14", \
RUN+="/bin/sh -c 'echo 180 > /sys$$DEVPATH/timeout'"
